# home/.config/systemd/user/homeshick-report.service
[Unit]
Description=Email a homeshick status report

[Service]
ExecStart=%h/.local/bin/homeshick-report
