# home/.config/systemd/user/in-work.target
[Unit]
Description=At the dayjob
Conflicts=in-work.timer
Wants=out-of-work.timer
