# home/.config/systemd/user/out-of-work.target
[Unit]
Description=Not at the dayjob
Conflicts=out-of-work.timer
Wants=in-work.timer
